# deployment/configs/supervisor/echo-server.conf
[program:echo-server]
directory=/var/www/app
command=laravel-echo-server start
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/echoserver.log
